# bdiofeed/bdiohid/files/etc/init.d/bdiohid
#!/bin/sh /etc/rc.common

START=56
STOP=54

start() {
    /usr/bin/bdiohid </dev/null &
}

stop() {
    /usr/bin/killall bdiohid
}