# bdiofeed/brltty/files/etc/init.d/brltty
#!/bin/sh /etc/rc.common

START=55
STOP=55

start() {
    /usr/bin/brltty -M1
}

stop() {
    /usr/bin/killall brltty
}